# templates/hugo.path.j2
[Unit]
Description=Watch the hugo filesystem for changes

[Path]
PathModified={{ hugo_dir }}

[Install]
WantedBy=multi-user.target
